# scripts/services/portforward/openclaw-portforward.service
[Unit]
Description=OpenClaw Dashboard Port-Forward via Tailscale
After=network.target k3s.service tailscaled.service
Wants=k3s.service tailscaled.service

[Service]
Type=simple
User=root
# Usar script wrapper para obtener IP de Tailscale dinámicamente
ExecStart=/opt/openclaw/portforward-wrapper.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
